# container/rc.local
#!/bin/sh -e
chown -R jenkins.jenkins /data

exit 0
