# host/builder-2.conf
description "builder 2"

start on started builder
stop on stopping builder
respawn

script
  /usr/bin/docker start -a builder-2
end script
